# configs/docker-images/devbox-nginx-php5.4/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php5-fpm
autostart=true
autorestart=true
priority=5

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=10

[program:ssh-d]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/cron -f -L 1
user=root
autostart=true
